# systemd/pipewire-web-remote.service
[Unit]
Description=PipeWire Web Remote Service
Documentation=https://github.com/oudeis01/pipewire-web-remote
After=pipewire.service
Requires=pipewire.service

[Service]
Type=simple
ExecStart=/usr/bin/pipewire-web-remote --listen 127.0.0.1:8449
Restart=on-failure
RestartSec=5

NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=default.target
